# etc/MAX-TEAM.conf
description "MAX-TEAM upstart script"

respawn
respawn limit 15 5

start on runlevel [2345]
stop on shutdown

setuid MAX-TEAM
exec /bin/sh MAX-TEAM/launch.sh
